# update-motd.d/10-help-text
#!/bin/sh

echo
echo "Welcome to Linux Mint"
echo " * Documentation:  http://www.linuxmint.com"

